# config/supervisord.conf

[supervisord]
nodaemon = true
user = root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
loglevel = INFO

[program:midserver]
command =/opt/agent/bin/wrapper-linux-x86-64 -c /opt/agent/conf/wrapper.conf wrapper.java.maxmemory=512 wrapper.syslog.ident=mid 
user = root
